# deploy/supervisor.conf
[program:ssp-backend-blue]
command=bash -c "cd /opt/ssp/backend && export $(openssl enc -aes-256-cbc -pbkdf2 -iter 100000 -d -in .env.enc -pass file:/etc/ssp/master.key | grep -v '^#' | xargs) && exec venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000"
directory=/opt/ssp/backend
user=ssp-app
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
stopwaitsecs=15
stderr_logfile=/var/log/ssp-backend-blue.err.log
stdout_logfile=/var/log/ssp-backend-blue.out.log

[program:ssp-backend-green]
command=bash -c "cd /opt/ssp/backend && export $(openssl enc -aes-256-cbc -pbkdf2 -iter 100000 -d -in .env.enc -pass file:/etc/ssp/master.key | grep -v '^#' | xargs) && exec venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001"
directory=/opt/ssp/backend
user=ssp-app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=15
stderr_logfile=/var/log/ssp-backend-green.err.log
stdout_logfile=/var/log/ssp-backend-green.out.log

[program:ssp-frontend-blue]
command=bash -c "PORT=3000 exec npm start"
directory=/opt/ssp/frontend
user=ssp-app
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
stopwaitsecs=15
stderr_logfile=/var/log/ssp-frontend-blue.err.log
stdout_logfile=/var/log/ssp-frontend-blue.out.log

[program:ssp-frontend-green]
command=bash -c "PORT=3002 exec npm start"
directory=/opt/ssp/frontend
user=ssp-app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=15
stderr_logfile=/var/log/ssp-frontend-green.err.log
stdout_logfile=/var/log/ssp-frontend-green.out.log
