# deploy/supervisor.conf
[program:ssp-backend-blue]
command=bash -c "cd /root/ssp/backend && export $(openssl enc -aes-256-cbc -pbkdf2 -iter 100000 -d -in .env.enc -pass file:/root/.ssp_master_key | grep -v '^#' | xargs) && exec venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000"
directory=/root/ssp/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/ssp-backend-blue.err.log
stdout_logfile=/var/log/ssp-backend-blue.out.log

[program:ssp-backend-green]
command=bash -c "cd /root/ssp/backend && export $(openssl enc -aes-256-cbc -pbkdf2 -iter 100000 -d -in .env.enc -pass file:/root/.ssp_master_key | grep -v '^#' | xargs) && exec venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001"
directory=/root/ssp/backend
autostart=false
autorestart=false
stderr_logfile=/var/log/ssp-backend-green.err.log
stdout_logfile=/var/log/ssp-backend-green.out.log

[program:ssp-frontend-blue]
command=bash -c "fuser -k 3000/tcp 2>/dev/null; sleep 1; PORT=3000 npm start"
directory=/root/ssp/frontend
autostart=true
autorestart=true
stderr_logfile=/var/log/ssp-frontend-blue.err.log
stdout_logfile=/var/log/ssp-frontend-blue.out.log

[program:ssp-frontend-green]
command=bash -c "fuser -k 3002/tcp 2>/dev/null; sleep 1; PORT=3002 npm start"
directory=/root/ssp/frontend
autostart=false
autorestart=false
stderr_logfile=/var/log/ssp-frontend-green.err.log
stdout_logfile=/var/log/ssp-frontend-green.out.log
