# deploy/uploads-gc.cron.example
# uploads 磁盘清理 cron(隐藏雷 #1)
#
# 安装:
#   crontab -l > /tmp/cron.bak 2>/dev/null || true
#   cat /tmp/cron.bak /opt/ssp/deploy/uploads-gc.cron.example | sort -u | crontab -
#   crontab -l                                # 验证
#   tail -f /var/log/ssp-uploads-gc.log       # 看输出
#
# 时间:每天 04:00(避开 03:00 备份窗口)
# 默认保留 90 天;改环境变量 SSP_UPLOADS_RETENTION_DAYS 调

0 4 * * * /usr/bin/bash /opt/ssp/deploy/uploads-gc.sh
